# SG-Router/VPS/wg-easy-home-route.timer.example
[Unit]
Description=Keep home LAN route through wg0

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
